# root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash
dnsmasq --hostsdir=${HOSTS_DIR} --no-daemon
