# root/etc/services.d/docker_monitor/run
#!/usr/bin/with-contenv bash
cd /app
python3 docker_monitor.py
